# systemd/.config/systemd/user/rclone-obsidian.service
[Unit]
Description=RClone Mount for Obsidian Vault
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=/bin/mkdir -p %h/.vaults/geistesblitze
ExecStart=/usr/bin/rclone mount gdrive-geistesblitze:Geistesblitze %h/.vaults/geistesblitze \
  --vfs-cache-mode writes \
  --allow-other \
  --allow-non-empty \
  --log-file %h/.logs/rclone-obsidian.log \
  --log-level INFO
ExecStop=/bin/fusermount -u %h/.vaults/geistesblitze
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
